# Docker/images/hadoop/config/supervisor/programs.conf
[program:start]
command=sh /root/start.sh
autostart=true
autorestart=false
startsecs = 0
nodaemon=true

[program:cron]
command=/usr/sbin/cron -f -L 15
autostart=true
autorestart=true

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true